# meetingbar.service
[Unit]
Description=MeetingBar – Calendar overlay for Hyprland + Waybar
# default.target is reached earlier and more reliably than graphical-session.target
After=default.target

[Service]
Type=simple
ExecStart={INSTALL_DIR}/.venv/bin/python3 {INSTALL_DIR}/daemon.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Pass Wayland/display env vars into the service (and to overlay subprocess).
# Make sure to add this to your hyprland.conf exec-once:
#   exec-once = systemctl --user import-environment WAYLAND_DISPLAY XDG_RUNTIME_DIR DISPLAY DBUS_SESSION_BUS_ADDRESS
PassEnvironment=WAYLAND_DISPLAY XDG_RUNTIME_DIR DISPLAY DBUS_SESSION_BUS_ADDRESS

[Install]
WantedBy=default.target
